# certbot.timer
[Unit]
Description=weekly check for letsencrypt renewals

[Timer]
OnCalendar=*-1/2-01 00:00:00
Persistent=True

[Install]
WantedBy=timers.target
